# todo/dunst/files/battery-check.timer
[Unit]
Description=Regularly check battery percentage

[Timer]
OnBootSec=3min
OnUnitActiveSec=3min
Persistent=true

[Install]
WantedBy=timers.target
